# services/memory-wiki/ai-chinese-server.service
[Unit]
Description=Unified Ai-Chinese Server (MemoryWiki API & Static Webapps)
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/alex/Ai-chinese/services/memory-wiki
ExecStart=/home/alex/.local/bin/uv run uvicorn app.main:app --host 0.0.0.0 --port 8082
Restart=on-failure

[Install]
WantedBy=default.target
